# amir2
#!/bin/sh /etc/rc.common

start(){

sh /root/owo.sh

}

reload(){

sh /root/up.sh

}
